# ersatztv/templates/ersatztv.service
[Unit]
Description=ErsatzTV Service
After=network.target

[Service]
ExecStart=~/ersatztv/ErsatzTV
Restart=always
StandardOutput=journal

[Install]
WantedBy=multi-user.target
